# systemd/system/pacman-F-sync.timer
[Unit]
Description=Weekly pacman -Fy

[Timer]
OnCalendar=weekly
AccuracySec=1d
Persistent=true

[Install]
WantedBy=timers.target
