# systemd/system/pacman-F-sync.service
[Unit]
Description=Run pacman -Fy
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Fy
