# systemd/system/powertop.timer
[Unit]
Description=Powertop tunings on boot

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
